# infra/systemd/knpodly.service
# Install to /etc/systemd/system/knpodly.service
# Then: systemctl daemon-reload && systemctl enable --now knpodly.service
[Unit]
Description=Knpodly Educational Linux Lab Platform (Docker Compose stack)
Requires=docker.service libvirtd.service
After=docker.service libvirtd.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/srv/knpodly
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml --env-file .env up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml --env-file .env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
